# packaging/app-file-scan.cron
0 2 * * * root /usr/sbin/file_scan >/dev/null 2>&1
